# mavlink-router.service.in
[Unit]
Description=MAVLink Router
Wants=mavlink_util.service

[Service]
Type=simple
ExecStart=@BINDIR@/mavlink-routerd -d /data/etc/companion/mavlink-router/
Restart=on-failure

[Install]
WantedBy=multi-user.target
